# app-metrics/keepalived_exporter_non-root/files/keepalived_exporter.initd
#!/sbin/openrc-run
# Copyright 2024 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="keepalive_exporter daemon"
description="Prometheus exporter for Keepalived metrics."
pidfile="/run/keepalived_exporter.pid"
command="/usr/bin/keepalived_exporter"
output_log="/var/log/keepalived_exporter/keepalived_exporter.log"
error_log="/var/log/keepalived_exporter/keepalived_exporter.err.log"
start_stop_daemon_args="--background --user keepalived_exporter --make-pidfile --pidfile ${pidfile}"

depend() {
    need net
}
