# app-admin/kolide-k2-bin/files/kolide-k2.initd
#!/sbin/openrc-run
# Copyright 2020-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="kolide-k2 daemon"

pidfile=/run/kolide-k2.pid
command="/opt/kolide-k2/bin/launcher"
command_args="-config /etc/kolide-k2/launcher.flags"
command_background="true"
output_log="/var/log/kolide-k2.log"
error_log="/var/log/kolide-k2.log"

depend() {
	after net
}
