# net-analyzer/zabbix_bot/files/init.d/zabbix_bot
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
    need localmount
}

start() {
    ebegin "Starting Zabbix Bot"
    start-stop-daemon --background --start --exec /usr/bin/zabbix_xmpp \
    --make-pidfile --pidfile /var/run/zabbix_xmpp
    eend $?
}

stop() {
    ebegin "Stoping Zabbix Bot"
    start-stop-daemon --stop --exec /usr/bin/zabbix_xmpp \
    --pidfile /var/run/zabbix_xmpp -s 1
    eend $?
}
