# net-analyzer/logdna-agent/files/logdna.init
#!/sbin/openrc-run

name="logdna-agent"
command="/usr/bin/${name}"
command_background='true'
pidfile="/var/run/${name}.pid"
logfile="/var/log/${name}.log"
output_log="${logfile}"
error_log="${logfile}"
