# app-admin/rasdaemon/files/ras-mc-ctl.init
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need rasdaemon
}

command=/usr/sbin/ras-mc-ctl
command_args=--register-labels
description='Initialize EDAC v3.0.0 Drivers For Machine Hardware'
in_background_fake=$command
